# backup/systemd_services/backup.timer
[Unit]
Description=Run backup

[Timer]
OnBootSec=5min
OnUnitActiveSec=6hours

[Install]
WantedBy=timers.target
